# config/includes.chroot/etc/systemd/system/lightdm.service.d/koppa-setup.conf
[Service]
ExecStartPre=/usr/local/bin/koppa-setup.sh
